# .config/systemd/user/padel.service
[Unit]
Description=Check for Padel court bookings

[Service]
Type=oneshot
ExecStart=padel

[Install]
WantedBy=multi-user.target
